# packaging/sing-box.service.d/10-sing-gateway.conf
[Service]
ExecStartPre=+/usr/bin/sing-gateway check
ExecStartPost=+/usr/bin/sing-gateway set
# Post-stop cleanup is best-effort so cleanup failure does not turn an
# otherwise clean sing-box stop into a failed unit state.
ExecStopPost=+-/usr/bin/sing-gateway unset
