# packaging/sing-box.service.d/10-sing-gateway.conf
[Service]
ExecStartPre=+/usr/bin/sing-gateway check
ExecStartPost=+/usr/bin/sing-gateway set
ExecStopPost=+/usr/bin/sing-gateway unset
